# static/crontab
# m h  dom mon dow   command

# Every hour
* * * * * python ping.py