# static/crontab
# m h  dom mon dow   command

# Every hour
*/2 * * * * python ping.py
0 */6 * * * sudo reboot -n

